# cloud/gateway-go/deployments/blockly-gateway.service
[Unit]
Description=Blockly Gateway Service (Golang)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/blockly-gateway
ExecStart=/opt/blockly-gateway/bin/gateway
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
